# packages/tools/volumectrl/system.d/volumectrl.service
[Unit]
Description=Volume control service

[Service]
Type=simple
ExecStart=/usr/bin/volumectrl

[Install]
WantedBy=multi-user.target
